# docker/all-in-one/services.d/mcp-server/run
#!/command/with-contenv sh
cd /app
exec python /app/mcp_entrypoint.py 2>&1
